# magic-conch.service
[Unit]
Description=Magic Conch Graveyard Museum
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/magic_conch
Environment="PATH=/home/YOUR_USERNAME/miniconda3/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/YOUR_USERNAME/miniconda3/bin/uvicorn server.app:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/var/log/magic-conch/access.log
StandardError=append:/var/log/magic-conch/error.log

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
